# deploy/celery_worker.service
[Unit]
Description=Voice AI Agent — Celery Worker
After=network.target redis.service postgresql.service

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/voice_ai_agent
Environment="DJANGO_SETTINGS_MODULE=config.settings.production"
ExecStart=/opt/voice_ai_agent/venv/bin/celery \
          -A config worker \
          --loglevel=info \
          --concurrency=4 \
          --logfile=/var/log/voice_ai_agent/celery_worker.log
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
