# deploy/voice_ai_gunicorn.service
[Unit]
Description=Voice AI Agent — Gunicorn WSGI server
Documentation=https://docs.gunicorn.org
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=agent
Group=agent
WorkingDirectory=/home/agent/voice_ai_agent
Environment="DJANGO_SETTINGS_MODULE=config.settings.production"

ExecStart=/home/agent/voice_ai_agent/venv/bin/gunicorn \
    --workers 3 \
    --bind 0.0.0.0:8000 \
    config.wsgi:application

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=voice_ai_gunicorn

[Install]
WantedBy=multi-user.target
